# deploy/jri.socket
[Unit]
Description=Just Ralph It socket

[Socket]
ListenStream=127.0.0.1:8000
NoDelay=true

[Install]
WantedBy=sockets.target
