# deploy/jri.service
[Unit]
Description=Just Ralph It
After=network.target

[Service]
Type=simple
User=nico
WorkingDirectory=/home/nico/jri
ExecStart=/usr/bin/python3 -m uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5
Environment=PATH=/home/nico/.local/bin:/home/linuxbrew/.linuxbrew/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=multi-user.target
